# config/supervisor.conf
[program:interviewquestions]
directory=/home/ben/interviewquestions
command=python manage.py run_gunicorn -c config/gunicorn.conf.py
user=root
autostart=true
autorestart=true
